# crontab
# Run at everyhours
0 * * * * /app/crawler >> /app/log.txt 2>&1
