# crontab
# Run at everyhours
* * * * * /app/crawler >> /app/log.txt 2>&1
