# crontab
# Run at 0h everyday
0 0 * * * /app/crawler >> /app/log.txt 2>&1
